# .config/autostart/background.desktop
[Desktop Entry]
Type=Application
Name=Background
Exec=bash -c 'feh --bg-scale ~/Pictures/Wallpaper.jpg'
