# .config/autostart/clipmenu.desktop
[Desktop Entry]
Type=Application
Name=Clipmenu Daemon
Exec=clipmenud
