# .config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Exec=bash -c 'conky'
