# .config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=Dunst daemon
Exec=dunst
