# .config/autostart/dwm-bar.desktop
[Desktop Entry]
Type=Application
Name=dwm-bar
Exec=bash -c '~/.config/dwm-bar/dwm_bar.sh' &
StartupWMClass="dwm"
