# .config/autostart/fluxgui.desktop
[Desktop Entry]
Icon=fluxgui
Type=Application
Name=f.lux indicator applet
Exec=fluxgui
