# .config/autostart/screenlayout.desktop
[Desktop Entry]
Type=Application
Name=Screenlayout
Exec=bash -c 'sleep 1 && ~/.screenlayout/default.sh'
