# .config/autostart/urxvtd.desktop
[Desktop Entry]
Type=Application
Name=Urxvt daemon
Exec=urxvtd
