# .config/autostart/windowwallpaper.desktop
[Desktop Entry]
Type=Application
Name=Background
Exec=st -t windowwallpaper -e bpytop
