# .config/autostart/xscreensaver.desktop
[Desktop Entry]
Type=Application
Name=Xscreensaver
Exec=xscreensaver -nosplash
